# components/build/supervisord.conf
[supervisord]
nodaemon=true
logfile = /var/run/supervisord.log
pidfile = /var/run/supervisord.pid

[program:sshd]  
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
startretries=10
exitcodes=0
stopsignal=KILL
stopwaitsecs=10
redirect_stderr=true

